# obd-autoconnect.service
[Unit]
Description=AutoPulse OBD Scanner Auto-Connect Service
After=bluetooth.service
Requires=bluetooth.service
Wants=bluetooth.service

[Service]
Type=oneshot
User=rocketeers
Group=rocketeers
WorkingDirectory=/home/rocketeers/vehicle_diagnostic_system
ExecStart=/home/rocketeers/vehicle_diagnostic_system/obd_autoconnect.sh connect
RemainAfterExit=yes
TimeoutStartSec=120

# Restart on failure
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obd-autoconnect

[Install]
WantedBy=multi-user.target
