# backend/vehicle-diagnostic.service
[Unit]
Description=Vehicle Diagnostic System - Auto Data Collection & ML Prediction
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=rocketeers
WorkingDirectory=/home/rocketeers/vehicle_diagnostic_system
Environment="PATH=/home/rocketeers/vehicle_diagnostic_system/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Start the complete system (backend servers + data collector)
ExecStart=/home/rocketeers/vehicle_diagnostic_system/auto_start_complete.sh

# Restart on failure
Restart=on-failure
RestartSec=10

# Graceful shutdown
ExecStop=/home/rocketeers/vehicle_diagnostic_system/stop_system.sh
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vehicle-diagnostic

[Install]
WantedBy=multi-user.target
